# netatalk/rootfs/etc/services.d/afp/run
#!/usr/bin/with-contenv bashio
exec netatalk -F /etc/afp.conf -d
