# victoriametrics/rootfs/etc/services.d/victoria-metrics/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run the Victoria Metrics
# ==============================================================================
exec victoria-metrics-prod --influxListenAddr=:8086 --httpListenAddr=:8428 -storageDataPath=/data -retentionPeriod="$(bashio::config 'retention')" -promscrape.config=/config/victoria-metrics/prometheus.yml
