# victoriametrics/rootfs/etc/services.d/victoria-metrics/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run the Victoria Metrics
# ==============================================================================
exec victoria-metrics-prod --influxListenAddr=:8089 --httpListenAddr=:8428 -storageDataPath=/data -retentionPeriod="$(bashio::services 'retention')"
